# ops/systemd/aplus-backfill.service
[Unit]
Description=Aplus OHLCV Backfill (Kraken via ccxt)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/aplus/aplus.env
WorkingDirectory=/srv/aplus
ExecStart=/usr/bin/python3 -m policy.backfill_ohlcv
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/aplus-backfill.log
StandardError=append:/var/log/aplus-backfill.log

[Install]
WantedBy=multi-user.target
