# ops/systemd/aplus-train.service
[Unit]
Description=Aplus Nightly Training
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/aplus/aplus.env
WorkingDirectory=/srv/aplus
Environment=TIMESTEPS=2000
ExecStart=/usr/bin/python3 -m policy.train
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/aplus-train.log
StandardError=append:/var/log/aplus-train.log

[Install]
WantedBy=multi-user.target
