# ops/systemd/dashboard.service
[Unit]
Description=TradePulse IQ FastAPI Dashboard
After=network.target

[Service]
ExecStart=/srv/trad/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
WorkingDirectory=/srv/trad
Restart=always
User=root
Group=root
EnvironmentFile=/etc/trad/trad.env
StandardOutput=append:/var/log/trad-dashboard.log
StandardError=append:/var/log/trad-dashboard.log

[Install]
WantedBy=multi-user.target
