# ops/systemd/trad-train.service
[Unit]
Description=Trad Nightly Training
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/trad/trad.env
WorkingDirectory=/srv/trad
Environment=TIMESTEPS=2000
ExecStart=/usr/bin/python3 -m policy.train
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/trad-train.log
StandardError=append:/var/log/trad-train.log

[Install]
WantedBy=multi-user.target
