# ops/systemd/trad-backfill.service
[Unit]
Description=Trad OHLCV Backfill (Kraken via ccxt)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/trad/trad.env
WorkingDirectory=/srv/trad
ExecStart=/usr/bin/python3 -m policy.backfill_ohlcv
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/trad-backfill.log
StandardError=append:/var/log/trad-backfill.log

[Install]
WantedBy=multi-user.target
