# ops/systemd/trad-meta.service
[Unit]
Description=Trad Meta-Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/trad/trad.env
WorkingDirectory=/srv/trad
ExecStart=/usr/bin/python3 -m meta_controller.run_metacontroller
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/trad-meta.log
StandardError=append:/var/log/trad-meta.log

[Install]
WantedBy=multi-user.target
