# ops/systemd/trad-auditor.service
[Unit]
Description=Run the Trad Pattern Auditor
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/srv/trad/.venv/bin/python /srv/trad/ops/scripts/run_auditor.py
WorkingDirectory=/srv/trad
Restart=on-failure
User=root
Group=root
EnvironmentFile=/etc/trad/trad.env

[Install]
WantedBy=multi-user.target
