# ops/systemd/trad-meta.timer
[Unit]
Description=Run Trad Meta-Controller hourly at :02

[Timer]
OnCalendar=*-*-* *:02:00
Unit=trad-meta.service

[Install]
WantedBy=timers.target
